# your_projects/mohdanwaroddin_WaterReminder/waterreminder.service
[Unit]
Description=Water Reminder service

[Service]
ExecStart=/usr/bin/WaterReminder.py
